# official_writeup/web-graphauth/game/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python /app/app.py

[program:auth]
command=python /app/auth.py
